# openbsd/lang302
#!/bin/ksh

rc_bg="YES"
daemon="/usr/local/bin/lang302"
lang302_user="_lang302"

. /etc/rc.d/rc.subr

rc_cmd $1
